# launchd/templates/jp.ontheroad.tmux-dashboard.tunnel.prod.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>jp.ontheroad.tmux-dashboard.tunnel.prod</string>

  <key>ProgramArguments</key>
  <array>
    <string>__REPO_ROOT__/launchd/start-tunnel-prod.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/launchd/logs/tunnel.prod.out.log</string>

  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/launchd/logs/tunnel.prod.err.log</string>
</dict>
</plist>
